# currency.crontab
@hourly curl -s 'http://api.fixer.io/latest' > /path/to/exchange.json
@hourly curl -s 'https://api.cryptonator.com/api/ticker/btc-eur' > /path/to/exchanges/btc.json
@hourly curl -s 'https://api.cryptonator.com/api/ticker/etc-eur' > /path/to/exchanges/etc.json
@hourly curl -s 'https://api.cryptonator.com/api/ticker/eth-eur' > /path/to/exchanges/eth.json
